# updating/my_overlay/etc/network/if-up.d/S01
#!/bin/sh
if [ "$IFACE" = eth0 ]; then
	ntpdate -u 0.pl.pool.ntp.org
	echo "ntp has been triggered"
fi
